# backend/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/nginx/error.log
stdout_logfile=/var/log/nginx/access.log
priority=10

[program:auth-service]
command=/usr/local/bin/auth-service
autostart=true
autorestart=true
startretries=10
startsecs=5
stderr_logfile=/var/log/supervisor/auth-service.err.log
stdout_logfile=/var/log/supervisor/auth-service.out.log
priority=20

[program:core-service]
command=/usr/local/bin/core-service
autostart=true
autorestart=true
startretries=10
startsecs=5
stderr_logfile=/var/log/supervisor/core-service.err.log
stdout_logfile=/var/log/supervisor/core-service.out.log
priority=30

[program:github-service]
command=/usr/local/bin/github-service
autostart=true
autorestart=true
startretries=10
startsecs=5
stderr_logfile=/var/log/supervisor/github-service.err.log
stdout_logfile=/var/log/supervisor/github-service.out.log
priority=40

[program:ai-service]
command=/usr/local/bin/ai-service
autostart=true
autorestart=true
startretries=10
startsecs=5
stderr_logfile=/var/log/supervisor/ai-service.err.log
stdout_logfile=/var/log/supervisor/ai-service.out.log
priority=50

[program:notification-service]
command=/usr/local/bin/notification-service
autostart=true
autorestart=true
startretries=10
startsecs=5
stderr_logfile=/var/log/supervisor/notification-service.err.log
stdout_logfile=/var/log/supervisor/notification-service.out.log
priority=60
