# backend/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/nginx/error.log
stdout_logfile=/var/log/nginx/access.log
priority=10

[program:auth-service]
command=/usr/local/bin/auth-service
autostart=true
autorestart=true
startretries=10
startsecs=5
stderr_logfile=/var/log/supervisor/auth-service.err.log
stdout_logfile=/var/log/supervisor/auth-service.out.log
priority=20
environment=PORT="8080",GITHUB_CLIENT_ID="%(ENV_GITHUB_CLIENT_ID)s",GITHUB_CLIENT_SECRET="%(ENV_GITHUB_CLIENT_SECRET)s",JWT_SECRET="%(ENV_JWT_SECRET)s",FRONTEND_URL="%(ENV_FRONTEND_URL)s",CORE_SERVICE_URL="%(ENV_CORE_SERVICE_URL)s"

[program:core-service]
command=/usr/local/bin/core-service
autostart=true
autorestart=true
startretries=10
startsecs=5
stderr_logfile=/var/log/supervisor/core-service.err.log
stdout_logfile=/var/log/supervisor/core-service.out.log
priority=30
environment=SERVER_PORT="8083",DATABASE_URL="%(ENV_DATABASE_URL)s",JWT_SECRET="%(ENV_JWT_SECRET)s",GITHUB_SERVICE_URL="%(ENV_GITHUB_SERVICE_URL)s",AI_SERVICE_URL="%(ENV_AI_SERVICE_URL)s",NOTIFICATION_SERVICE_URL="%(ENV_NOTIFICATION_SERVICE_URL)s"

[program:github-service]
command=/usr/local/bin/github-service
autostart=true
autorestart=true
startretries=10
startsecs=5
stderr_logfile=/var/log/supervisor/github-service.err.log
stdout_logfile=/var/log/supervisor/github-service.out.log
priority=40
environment=PORT="8081",GITHUB_TOKEN="%(ENV_GITHUB_TOKEN)s"

[program:ai-service]
command=/usr/local/bin/ai-service
autostart=true
autorestart=true
startretries=10
startsecs=5
stderr_logfile=/var/log/supervisor/ai-service.err.log
stdout_logfile=/var/log/supervisor/ai-service.out.log
priority=50
environment=AI_SERVICE_PORT="8082",GEMINI_API_KEY="%(ENV_GEMINI_API_KEY)s",GROQ_API_KEY="%(ENV_GROQ_API_KEY)s",LLM_PROVIDER="%(ENV_LLM_PROVIDER)s",LLM_MODEL="%(ENV_LLM_MODEL)s"

[program:notification-service]
command=/usr/local/bin/notification-service
autostart=true
autorestart=true
startretries=10
startsecs=5
stderr_logfile=/var/log/supervisor/notification-service.err.log
stdout_logfile=/var/log/supervisor/notification-service.out.log
priority=60
environment=PORT="8084"
